# g935-dsp.service
[Unit]
Description=G935 DSP soundstage + mic daemon (HID++ power-on watcher)

[Service]
# Prefer `make install-user` / `./install.sh --user`, which writes a unit with
# the absolute path to the installed wrapper. This checked-in unit assumes the
# wrapper is on PATH (typical after install to ~/.local/bin).
ExecStart=%h/.local/bin/g935-dspd
Restart=always
RestartSec=5
SyslogIdentifier=g935-dspd

[Install]
WantedBy=default.target
